# nightly/turn-off-the-lights.timer
[Unit]
Description=turn off the V0 testbed on a daily basis

[Timer]
OnCalendar=20:00

[Install]
WantedBy=multi-user.target
