# services/pull-and-restart.service
[Unit]
Description="Pull from git and restart services"

[Service]
ExecStart=/root/r2lab-embedded/services/pull-and-restart.sh
