# services/pull-and-restart.timer
[Unit]
Description="Pull from git and restart services"

[Timer]
Unit=pull-and-restart.service
OnCalendar=06:00

[Install]
WantedBy=timers.target
